# templates/rpm-sources/bin/init.d-mesos-slaved.sh.in
#!/bin/bash

# Licensed to the Apache Software Foundation (ASF) under one
# or more contributor license agreements.  See the NOTICE file
# distributed with this work for additional information
# regarding copyright ownership.  The ASF licenses this file
# to you under the Apache License, Version 2.0 (the
# "License"); you may not use this file except in compliance
# with the License.  You may obtain a copy of the License at
# 
#     http://www.apache.org/licenses/LICENSE-2.0
# 
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

# This is a wrapper for running mesos-slave before it is installed
# that first sets up environment variables as appropriate.


# mesos-slaved - Startup script for Mesos Slave

# chkconfig: 35 85 15
# description: 
# processname: mesos-slaved 
# config: /etc/mesos.conf
# pidfile: /var/run/mesos/mesos-slaved.pid

. /etc/rc.d/init.d/functions

# NOTE: if you change any OPTIONS here, you get what you pay for:
# this script assumes all options are in the config file.

ENVIRONMENT="/etc/profile.d/mesos"
if [ -x "$ENVIRONMENT" ]; then
	. "$ENVIRONMENT"
fi

CONFIGFILE="/etc/mesos/conf"
OPTIONS=" -f $CONFIGFILE"

mesosd=${MESOS_SLAVED-/usr/bin/mesos-slaved}

MESOS_USER=${MESOS_USER-mesosd}
MESOS_GROUP=${MESOS_GROUP-mesosd}

start()
{
    echo -n $"Starting Mesos Slave (mesos-slaved):"
    daemon --user "$MESOS_USER" $NUMACTL $mesosd $OPTIONS
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && touch /var/lock/subsys/mesos-slaved
}

stop()
{
    echo -n $"Stopping Mesos Slave (mesos-slaved): "
    killproc -p "$PIDFILE" -d 300 $mesosd
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/mesos-slaved
}

restart () {
	stop
	start
}

ulimit -n 12000
RETVAL=0

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart|reload|force-reload)
    restart
    ;;
  condrestart)
    [ -f /var/lock/subsys/mesos-slaved ] && restart || :
    ;;
  status)
    status $mesosd
    RETVAL=$?
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart|reload|force-reload|condrestart}"
    RETVAL=1
esac

exit $RETVAL
